# supervisor/templates/default/supervisor_service.conf.erb
[program:<%= @name %>]
command=<%= @command %>
environment=<%= @environment %>
